# docker_configs/crontab
*/5 * * * * /var/www/bin/concole nginx:reload